# iso-builder/nyx-profile/airootfs/etc/skel/.config/systemd/user/nyxus-missiond.service
[Unit]
Description=NYXUS Mission Control daemon
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/nyxus/nyxus_missiond.py
Restart=on-failure
RestartSec=2
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ReadWritePaths=%t/nyxus-mission %h/.config/nyxus %h/.cache/nyxus

[Install]
WantedBy=graphical-session.target
