# iso-builder/nyx-profile/airootfs/etc/skel/.config/systemd/user/nyxus-qsd.service
[Unit]
Description=NYXUS Quick Settings daemon
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/nyxus/nyxus_qsd.py
Restart=on-failure
RestartSec=2
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ReadWritePaths=%t/nyxus-qs %h/.config/nyxus

[Install]
WantedBy=graphical-session.target
