# artifacts/api-server/nyxus-scripts/nyxus-eww.service
[Unit]
Description=NYXUS EWW daemon (ElKowar's Wacky Widgets)
Documentation=https://github.com/elkowar/eww
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=simple
EnvironmentFile=-%h/.config/eww/nyxus.conf
ExecStartPre=/usr/bin/sh -c 'mkdir -p "${NYXUS_EWW_LOG_DIR:-$HOME/.cache/nyxus-eww}"'
ExecStart=/usr/bin/env eww daemon --no-daemonize
ExecStartPost=/usr/local/bin/nyxus-eww-launch
ExecStop=/usr/bin/env eww kill
Restart=on-failure
RestartSec=3
StandardOutput=append:%h/.cache/nyxus-eww/daemon.log
StandardError=append:%h/.cache/nyxus-eww/daemon.log

[Install]
WantedBy=graphical-session.target
