# artifacts/api-server/nyxus-scripts/nyxus-dynamic-wallpaper.timer
[Unit]
Description=Re-pick NYXUS dynamic wallpaper every hour

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h
OnCalendar=*:00
Persistent=true
Unit=nyxus-dynamic-wallpaper.service

[Install]
WantedBy=timers.target
